# production/app/supervisord.ini
[supervisord]
nodaemon=true
user=root

[program:nginx]
command=nginx
stdout_logfile=/var/www/html/storage/logs/nginx/stdout.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/www/html/storage/logs/nginx/stderr.log
stderr_logfile_maxbytes=0

[program:php-fpm]
command=php-fpm8
stdout_logfile=/var/www/html/storage/logs/php/stdout.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/www/html/storage/logs/php/stderr.log
stderr_logfile_maxbytes=0

[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work redis --max-time=3600 --timeout=20
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=32
stdout_logfile=/var/www/html/storage/logs/laravel-worker/stdout.log
stdout_logfile_maxbytes=0
stderr_logfile=/var/www/html/storage/logs/laravel-worker/stderr.log
stderr_logfile_maxbytes=0
stopwaitsecs=3600
